# systemd/sentry-worker.service
[Unit]
Description=Sentry Background Worker
After=network.target

[Service]
Type=simple
User=sentry
Group=daemon
WorkingDirectory=/var/opt/sentry
EnvironmentFile=/etc/default/sentry
ExecStart=/usr/sbin/sentry-worker

[Install]
WantedBy=multi-user.target
